# slash/etc/profile.d/export.sh
#<XDG>
export XDG_DATA_DIRS=$HOME/.nix-profile/share:/usr/local/share:/usr/share:/var/lib/snapd/desktop
#<JAVA>
export _JAVA_AWT_WM_NONREPARENTING=1
export JAVA_HOME=/usr/lib/jvm/jdk-19
